# config/supervisor.conf
[program:%(pasta)s]
command=%(virtualenv)s/bin/gunicorn twitsearch.wsgi -t120 -b 127.0.0.1:%(porta)s -w1 --pythonpath=%(fullpath)s
user=webapp
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
redirect_stderr=true